# scripts/systemd/lotto.timer
[Unit]
Description=Lotto Auto Purchase Timer - Every Sunday 09:00 AM KST
Requires=lotto.service

[Timer]
# Every Saturday at 09:00 AM KST
OnCalendar=Sat *-*-* 09:00:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
